# test/synthetic-gui-case/case.json
{
  "case_id": "case-synthetic-07145bf0967a4a28a52b542f407f540e",
  "case_name": "synthetic_behavior_case",
  "sample_name": "synthetic-sample.exe",
  "sample_path": "synthetic://no-executable",
  "sha256": "synthetic",
  "started_at": "2026-07-13T11:01:11.8329055+09:00",
  "duration": "00:00:02",
  "static_analysis": null,
  "processes": [
    {
      "name": "synthetic-sample.exe",
      "pid": 4242,
      "parent_pid": null,
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "command_line": "synthetic-sample.exe",
      "image_path": "synthetic://no-executable",
      "start_time": "2026-07-13T11:01:11.8329055+09:00",
      "end_time": "2026-07-13T11:01:13.8329055+09:00",
      "event_count": 14,
      "network_count": 0,
      "file_count": 2,
      "registry_count": 1
    }
  ],
  "events": [
    {
      "id": 1,
      "time": "00:00:00.1000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Process",
      "action": "Process Start",
      "object_value": "synthetic-sample.exe",
      "summary": "Synthetic process start",
      "severity": "Informational",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Process\u0022,\r\n  \u0022action\u0022: \u0022Process Start\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022synthetic-sample.exe\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic process start\u0022,\r\n  \u0022technique_id\u0022: \u0022\u0022,\r\n  \u0022technique_name\u0022: \u0022\u0022,\r\n  \u0022confidence\u0022: \u0022\u0022\r\n}",
      "technique_id": "",
      "technique_name": "",
      "confidence": "",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 2,
      "time": "00:00:00.3000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Api",
      "action": "VirtualAllocEx",
      "object_value": "pid-5000",
      "summary": "Synthetic remote allocation event",
      "severity": "High",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Api\u0022,\r\n  \u0022action\u0022: \u0022VirtualAllocEx\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022pid-5000\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic remote allocation event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1055\u0022,\r\n  \u0022technique_name\u0022: \u0022Process Injection\u0022,\r\n  \u0022confidence\u0022: \u0022High\u0022\r\n}",
      "technique_id": "T1055",
      "technique_name": "Process Injection",
      "confidence": "High",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 3,
      "time": "00:00:00.4000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Api",
      "action": "WriteProcessMemory",
      "object_value": "pid-5000",
      "summary": "Synthetic remote memory write event",
      "severity": "High",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Api\u0022,\r\n  \u0022action\u0022: \u0022WriteProcessMemory\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022pid-5000\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic remote memory write event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1055\u0022,\r\n  \u0022technique_name\u0022: \u0022Process Injection\u0022,\r\n  \u0022confidence\u0022: \u0022High\u0022\r\n}",
      "technique_id": "T1055",
      "technique_name": "Process Injection",
      "confidence": "High",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 4,
      "time": "00:00:00.5000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Api",
      "action": "CreateRemoteThread",
      "object_value": "pid-5000",
      "summary": "Synthetic remote execution event",
      "severity": "High",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Api\u0022,\r\n  \u0022action\u0022: \u0022CreateRemoteThread\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022pid-5000\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic remote execution event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1055\u0022,\r\n  \u0022technique_name\u0022: \u0022Process Injection\u0022,\r\n  \u0022confidence\u0022: \u0022High\u0022\r\n}",
      "technique_id": "T1055",
      "technique_name": "Process Injection",
      "confidence": "High",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 5,
      "time": "00:00:00.5010000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Behavior",
      "action": "Remote Thread Injection",
      "object_value": "T1055",
      "summary": "Remote allocation, process-memory write, and remote thread creation were observed. [rule_version=1; tactic=; rule_hash=af386374ea8cde1506505fdb0c3b615b9b7b363705411c770d1fb9e56d7d3579]",
      "severity": "High",
      "source": "Behavior",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Behavior\u0022,\r\n  \u0022action\u0022: \u0022Remote Thread Injection\u0022,\r\n  \u0022technique_id\u0022: \u0022T1055\u0022,\r\n  \u0022technique_name\u0022: \u0022Process Injection\u0022,\r\n  \u0022confidence\u0022: \u0022High\u0022,\r\n  \u0022evidence_event_ids\u0022: [\r\n    2,\r\n    3,\r\n    4\r\n  ],\r\n  \u0022evidence_actions\u0022: [\r\n    \u0022VirtualAllocEx\u0022,\r\n    \u0022WriteProcessMemory\u0022,\r\n    \u0022CreateRemoteThread\u0022\r\n  ]\r\n}",
      "technique_id": "T1055",
      "technique_name": "Process Injection",
      "confidence": "High",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 6,
      "time": "00:00:01",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "File",
      "action": "CreateFileW",
      "object_value": "C:\\Temp\\sample.doc.locked",
      "summary": "Synthetic file create",
      "severity": "Medium",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022File\u0022,\r\n  \u0022action\u0022: \u0022CreateFileW\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022C:\\\\Temp\\\\sample.doc.locked\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic file create\u0022,\r\n  \u0022technique_id\u0022: \u0022\u0022,\r\n  \u0022technique_name\u0022: \u0022\u0022,\r\n  \u0022confidence\u0022: \u0022\u0022\r\n}",
      "technique_id": "",
      "technique_name": "",
      "confidence": "",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 7,
      "time": "00:00:01.1000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "File",
      "action": "WriteFile",
      "object_value": "C:\\Temp\\sample.doc.locked",
      "summary": "Synthetic file write",
      "severity": "Medium",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022File\u0022,\r\n  \u0022action\u0022: \u0022WriteFile\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022C:\\\\Temp\\\\sample.doc.locked\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic file write\u0022,\r\n  \u0022technique_id\u0022: \u0022\u0022,\r\n  \u0022technique_name\u0022: \u0022\u0022,\r\n  \u0022confidence\u0022: \u0022\u0022\r\n}",
      "technique_id": "",
      "technique_name": "",
      "confidence": "",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 8,
      "time": "00:00:01.2000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Api",
      "action": "CryptEncrypt",
      "object_value": "buffer",
      "summary": "Synthetic crypto API event",
      "severity": "Medium",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Api\u0022,\r\n  \u0022action\u0022: \u0022CryptEncrypt\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022buffer\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic crypto API event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1486\u0022,\r\n  \u0022technique_name\u0022: \u0022Data Encrypted for Impact\u0022,\r\n  \u0022confidence\u0022: \u0022Low\u0022\r\n}",
      "technique_id": "T1486",
      "technique_name": "Data Encrypted for Impact",
      "confidence": "Low",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 9,
      "time": "00:00:01.4000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Registry",
      "action": "RegSetValueExW",
      "object_value": "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Synthetic",
      "summary": "Synthetic persistence registry event",
      "severity": "High",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Registry\u0022,\r\n  \u0022action\u0022: \u0022RegSetValueExW\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022HKCU\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\Synthetic\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic persistence registry event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1547\u0022,\r\n  \u0022technique_name\u0022: \u0022Boot or Logon Autostart Execution\u0022,\r\n  \u0022confidence\u0022: \u0022Medium\u0022\r\n}",
      "technique_id": "T1547",
      "technique_name": "Boot or Logon Autostart Execution",
      "confidence": "Medium",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 10,
      "time": "00:00:01.4010000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Behavior",
      "action": "Persistence Established",
      "object_value": "T1547",
      "summary": "A registry or service change associated with autostart persistence was observed.",
      "severity": "High",
      "source": "Behavior",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Behavior\u0022,\r\n  \u0022action\u0022: \u0022Persistence Established\u0022,\r\n  \u0022technique_id\u0022: \u0022T1547\u0022,\r\n  \u0022technique_name\u0022: \u0022Boot or Logon Autostart Execution\u0022,\r\n  \u0022confidence\u0022: \u0022Medium\u0022,\r\n  \u0022evidence_event_ids\u0022: [\r\n    8\r\n  ],\r\n  \u0022evidence_actions\u0022: [\r\n    \u0022RegSetValueExW\u0022\r\n  ]\r\n}",
      "technique_id": "T1547",
      "technique_name": "Boot or Logon Autostart Execution",
      "confidence": "Medium",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 11,
      "time": "00:00:01.6000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Api",
      "action": "AmsiScanBuffer",
      "object_value": "synthetic-script",
      "summary": "Synthetic AMSI event",
      "severity": "Medium",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Api\u0022,\r\n  \u0022action\u0022: \u0022AmsiScanBuffer\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022synthetic-script\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic AMSI event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1562.001\u0022,\r\n  \u0022technique_name\u0022: \u0022Disable or Modify Tools\u0022,\r\n  \u0022confidence\u0022: \u0022Medium\u0022\r\n}",
      "technique_id": "T1562.001",
      "technique_name": "Disable or Modify Tools",
      "confidence": "Medium",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 12,
      "time": "00:00:01.7000000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Api",
      "action": "VirtualProtect",
      "object_value": "memory page",
      "summary": "Synthetic executable memory permission event",
      "severity": "Medium",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Api\u0022,\r\n  \u0022action\u0022: \u0022VirtualProtect\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022memory page\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic executable memory permission event\u0022,\r\n  \u0022technique_id\u0022: \u0022T1562.001\u0022,\r\n  \u0022technique_name\u0022: \u0022Disable or Modify Tools\u0022,\r\n  \u0022confidence\u0022: \u0022Medium\u0022\r\n}",
      "technique_id": "T1562.001",
      "technique_name": "Disable or Modify Tools",
      "confidence": "Medium",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 13,
      "time": "00:00:02",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Process",
      "action": "Process Exit",
      "object_value": "synthetic-sample.exe",
      "summary": "Synthetic process exit",
      "severity": "Informational",
      "source": "Synthetic",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Synthetic\u0022,\r\n  \u0022synthetic\u0022: true,\r\n  \u0022category\u0022: \u0022Process\u0022,\r\n  \u0022action\u0022: \u0022Process Exit\u0022,\r\n  \u0022pid\u0022: 4242,\r\n  \u0022obj\u0022: \u0022synthetic-sample.exe\u0022,\r\n  \u0022summary\u0022: \u0022Synthetic process exit\u0022,\r\n  \u0022technique_id\u0022: \u0022\u0022,\r\n  \u0022technique_name\u0022: \u0022\u0022,\r\n  \u0022confidence\u0022: \u0022\u0022\r\n}",
      "technique_id": "",
      "technique_name": "",
      "confidence": "",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    },
    {
      "id": 14,
      "time": "00:00:02.0010000",
      "process": "synthetic-sample.exe",
      "pid": 4242,
      "category": "Behavior",
      "action": "Process Injection Detected",
      "object_value": "T1055",
      "summary": "Remote memory allocation, process-memory write, and remote execution APIs were observed in sequence.",
      "severity": "High",
      "source": "Behavior",
      "raw_json": "{\r\n  \u0022source\u0022: \u0022Behavior\u0022,\r\n  \u0022action\u0022: \u0022Process Injection Detected\u0022,\r\n  \u0022technique_id\u0022: \u0022T1055\u0022,\r\n  \u0022technique_name\u0022: \u0022Process Injection\u0022,\r\n  \u0022confidence\u0022: \u0022High\u0022,\r\n  \u0022evidence_event_ids\u0022: [\r\n    1,\r\n    2,\r\n    3,\r\n    4,\r\n    5,\r\n    6,\r\n    7,\r\n    8,\r\n    9,\r\n    10,\r\n    11\r\n  ],\r\n  \u0022evidence_actions\u0022: [\r\n    \u0022Process Start\u0022,\r\n    \u0022VirtualAllocEx\u0022,\r\n    \u0022WriteProcessMemory\u0022,\r\n    \u0022CreateRemoteThread\u0022,\r\n    \u0022CreateFileW\u0022,\r\n    \u0022WriteFile\u0022,\r\n    \u0022CryptEncrypt\u0022,\r\n    \u0022RegSetValueExW\u0022,\r\n    \u0022AmsiScanBuffer\u0022,\r\n    \u0022VirtualProtect\u0022,\r\n    \u0022Process Exit\u0022\r\n  ]\r\n}",
      "technique_id": "T1055",
      "technique_name": "Process Injection",
      "confidence": "High",
      "process_guid": "case-synthetic-07145bf0967a4a28a52b542f407f540e:4242:639195048718329055",
      "captured_at_utc": null
    }
  ],
  "artifacts": [
    {
      "type": "Api",
      "value": "pid-5000, buffer, synthetic-script, memory page",
      "first_seen": "2026-07-13T11:01:12.1329055+09:00",
      "last_seen": "2026-07-13T11:01:13.5329055+09:00",
      "event_count": 6,
      "related_processes": "synthetic-sample.exe",
      "severity": "High"
    },
    {
      "type": "Behavior",
      "value": "T1055, T1547",
      "first_seen": "2026-07-13T11:01:12.3339055+09:00",
      "last_seen": "2026-07-13T11:01:13.8339055+09:00",
      "event_count": 3,
      "related_processes": "synthetic-sample.exe",
      "severity": "High"
    },
    {
      "type": "File",
      "value": "C:\\Temp\\sample.doc.locked",
      "first_seen": "2026-07-13T11:01:12.8329055+09:00",
      "last_seen": "2026-07-13T11:01:12.9329055+09:00",
      "event_count": 2,
      "related_processes": "synthetic-sample.exe",
      "severity": "Medium"
    },
    {
      "type": "Registry",
      "value": "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Synthetic",
      "first_seen": "2026-07-13T11:01:13.2329055+09:00",
      "last_seen": "2026-07-13T11:01:13.2329055+09:00",
      "event_count": 1,
      "related_processes": "synthetic-sample.exe",
      "severity": "High"
    }
  ],
  "network_sessions": [],
  "analyst_notes": "Synthetic case only. No operating-system behaviors were executed.",
  "quality": null,
  "raw_evidence_files": null,
  "preserved_files": null,
  "raw_evidence": []
}